# ansible/template/supervisor.conf.j2
[program:{{ project_fullname }}]
command={{ script }}
directory={{ dest }}
