# ansible/template/supervisor.conf.j2
[program:{{ project_name }}]
command={{ venv_gunicorn }} -c {{ gunicorn_conf }} {{ wsgi_app }}
directory={{ dest }}
